# chromium.service
[Unit]
Description=chromium
Type=simple

[Service]
ExecStart=/usr/bin/chromium

[Install]
WantedBy=default.target

