# chromium.service
[Unit]
Description=chromium
Type=simple

[Service]
ExecStart=/usr/bin/chromium --force-device-scale-factor=1.2

[Install]
WantedBy=default.target

