# xautolock.service
[Unit]
Description=xautolock
Type=simple

[Service]
ExecStart=/usr/bin/xautolock -time 15 -locker slock

[Install]
WantedBy=default.target

